# systemd/user/hypr-monitor-resume.service
[Unit]
Description=Fix HDMI monitor after suspend for Hyprland
After=suspend.target hibernate.target hybrid-sleep.target suspend-then-hibernate.target

[Service]
Type=oneshot
ExecStart=%h/.config/hypr/scripts/fix-monitors.sh post

[Install]
WantedBy=suspend.target hibernate.target hybrid-sleep.target suspend-then-hibernate.target
